# wordpress/rootfs/etc/services.d/nginx/run
#!/usr/bin/env bashio

# Start NGiNX
bashio::log.info "Starting NGinx..."

exec nginx -g 'daemon off; pid /tmp/nginx.pid;'